# Operation System/Praktikum1/task-2/crontabs
0 * * * * chmod +x yu_database.sh && ./yu_database.sh
